# fs_overlay/etc/update-motd.d/00-header
#!/bin/sh
#
# 00-header - create the header of the MOTD
# Copyright appnostic.  All rights reserved.

[ -r /etc/os-release ] && . /etc/os-release

cat <<_EOF
███╗   ██╗ ██████╗ ███████╗████████╗██╗ ██████╗ ██████╗ ███████╗
████╗  ██║██╔═══██╗██╔════╝╚══██╔══╝██║██╔════╝██╔═══██╗██╔════╝
██╔██╗ ██║██║   ██║███████╗   ██║   ██║██║     ██║   ██║███████╗
██║╚██╗██║██║   ██║╚════██║   ██║   ██║██║     ██║   ██║╚════██║
██║ ╚████║╚██████╔╝███████║   ██║   ██║╚██████╗╚██████╔╝███████║
╚═╝  ╚═══╝ ╚═════╝ ╚══════╝   ╚═╝   ╚═╝ ╚═════╝ ╚═════╝ ╚══════╝
_EOF
printf "Welcome to %s %s (%s %s %s)\n" "$NAME" "$VERSION" "$(uname -o)" "$(uname -r)" "$(uname -m)"                                                            
